# cygwin/etc/profile.d/lv2.csh
if ( ! $?LV2_PATH ) setenv LV2_PATH ""
setenv LV2_PATH "/usr/lib/lv2"
